# deploy_tools/gunicorn-upstart-template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid calimat
chdir /home/tourngenserver/sites/SITENAME/source

exec ../virtualenv/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    tourngen.wsgi:application
